# geoport/supervisord.conf
[program:celery]
command=celery -A geoport worker -l info --workdir {{ PROJECT_DIR }}
